# scripts/watchdog.service
[Unit]
Description=Chat Anon Frontend + Backend Watchdog
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/admin/chat_anon_ws/chat_anon_backend
ExecStart=/home/admin/chat_anon_ws/chat_anon_backend/scripts/start-watchdog.sh
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
